# RAG_MCP/systemd/tab-mcp-client.service
[Unit]
Description=TA_V8 TAB MCP Client (Web UI)
After=network.target postgresql.service neo4j.service qdrant.service minio.service
Wants=postgresql.service neo4j.service qdrant.service minio.service
PartOf=ta-v8-rag.target

[Service]
Type=simple
User=mentorius
Group=mentorius
WorkingDirectory=/home/mentorius/AI_Services/TA_V8/RAG_MCP/TAB_MCP_Client
Environment="PYTHONPATH=/home/mentorius/AI_Services/TA_V8/RAG_MCP"
Environment="PATH=/home/mentorius/AI_Services/TA_V8/RAG_MCP/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/home/mentorius/AI_Services/TA_V8/RAG_MCP/.venv"

# Database and Storage Configuration  
Environment="POSTGRES_HOST=localhost"
Environment="POSTGRES_PORT=5432"
Environment="POSTGRES_DB=ta_v8"
Environment="POSTGRES_USER=postgres_user"
Environment="POSTGRES_PASSWORD=postgres_pass"
Environment="NEO4J_URI=bolt://localhost:7687"
Environment="NEO4J_USER=neo4j"
Environment="NEO4J_PASSWORD=pJnssz3khcLtn6T"
Environment="MINIO_ENDPOINT=localhost:9000"
Environment="MINIO_ACCESS_KEY=minioadmin"
Environment="MINIO_SECRET_KEY=minioadmin"
Environment="QDRANT_HOST=localhost"
Environment="QDRANT_PORT=6333"

# MCP Services Configuration - using stdio transport
Environment="MCP_TRANSPORT=stdio"
Environment="CHUNKING_MCP_SERVICE=chunking-mcp"
Environment="EMBEDDING_MCP_SERVICE=embedding-mcp"
Environment="RETRIEVAL_MCP_SERVICE=retrieval-mcp"
Environment="RAG_AGENT_TEAM_SERVICE=rag-agent-team"

# Web UI Configuration
Environment="API_SECRET_KEY=ta_v8_secure_api_key_prod_2024_xK9mP2vB8nQ4"

ExecStart=/home/mentorius/AI_Services/TA_V8/RAG_MCP/.venv/bin/python -m uvicorn backend.app:app --host 0.0.0.0 --port 8005
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tab-mcp-client

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=100%

[Install]
WantedBy=ta-v8-rag.target
